# NVIDIA/TD-9-11-22/Rendu Younousse/databasesaving.timer
[Unit]
Description=lance  dataabasesaving.service chaque 10min

[Timer]
OnBootSec=0
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
